# files/rc.local
#!/bin/sh -e

/etc/init.d/mysql start

exit 0